# build-dashboard-api.service
[Unit]
Description=Build Dashboard API (FastAPI)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/build/api
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8001
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=build-dashboard-api

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target

